# config/app.conf
[program:game]
command:/usr/bin/python -m game_server
directory: /home/ubuntu
autostart:true
